# voice-dictation.service
[Unit]
Description=Voice Dictation
After=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h/src/voice-dictation
ExecStart=%h/src/voice-dictation/run-daemon.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
